# ops/launchagents/com.amerikagundemi.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.amerikagundemi.pipeline</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/tugrultasgetiren/bin/amerikagundemi-pipeline</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Minute</key>
    <integer>5</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Volumes/Crucial X10/Claude Code Projects/AmerikaGundemi</string>
  <key>StandardOutPath</key>
  <string>/Users/tugrultasgetiren/.amerikagundemi/logs/launchagent-pipeline.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/tugrultasgetiren/.amerikagundemi/logs/launchagent-pipeline.err.log</string>
</dict>
</plist>
